# telegram-bot.service
[Unit]
Description=Telegram Bet Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/telegram_bet_bot
Environment="PATH=/home/YOUR_USERNAME/telegram_bet_bot/venv/bin"
ExecStart=/home/YOUR_USERNAME/telegram_bet_bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
